# myserver.service
[Unit]
Description=myserver
After=network.target

[Service]
ExecStart=/usr/local/bin/server
Restart=on-failure

[Install]
WantedBy=multi-user.target
